# keto/lib/systemd/system/ketod.service
[Unit]
Description=ketod
[Service]
ExecStart=/opt/keto/bin/ketod.sh
PrivateTmp=true
